# deb/postinst
#!/bin/bash

systemctl enable pastcodd.service
systemctl start pastcodd.service
