# deb/postinst
#!/bin/bash

update-rc.d pastcodd defaults 20 20
