# deb/prerm
systemctl stop pastcodd.service
systemctl disable pastcodd.service
